# docker/supervisor/supervisor.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/php-fpm-error.log
stdout_logfile=/var/log/supervisor/php-fpm.log

[program:rabbitmq-consumer]
command=php /var/www/html/bin/console rabbitmq:consumer update_statistics -m 10 -vvv
user=www-data
autostart=true
autorestart=true
numprocs=1
process_name=%(program_name)s
stderr_logfile=/var/log/supervisor/rabbitmq-error.log
stdout_logfile=/var/log/supervisor/rabbitmq.log
